# Docker/api.conf
[program:api]
user=hygieia
directory=/var/Hygieia/api
command=mvn jetty:run -DPROP_FILE=application.properties
autostart=true
autorestart=true
redirect_stderr=true
