# Docker/jenkins.conf
[program:jenkins]
user=hygieia
directory=/var/Hygieia/BuildCollector/target/
command=java -jar hudson-collector-0.1.1.jar
autostart=true
autorestart=true
redirect_stderr=true
