# Docker/mongodb.conf
[program:mongod]
directory=/data/db
command=mongod
autostart=true
autorestart=true
redirect_stderr=true
