# Docker/sonar.conf
[program:sonar]
user=hygieia
directory=/var/Hygieia/CodeQualityCollector/target/
command=java -jar sonar-collector-0.1.0.jar
autostart=true
autorestart=true
redirect_stderr=true
