# Docker/ui.conf
[program:ui]
user=hygieia
environment=HOME="/home/hygieia",USER="hygieia"
directory=/var/Hygieia/UI
command=gulp serve
autostart=true
autorestart=true
redirect_stderr=true
